# templates/default/cloudwatch_logs/awslogs.init.sh.erb
#!/bin/sh
#
# awslogs               Start/Stop the AWSLogs agent
#
# chkconfig: 2345 95 25
# description:      This script is responsible for running AWSLogs agent \
#                   as a daemon
### BEGIN INIT INFO
# Provides:          awslogs
# Required-Start:    networking
# Required-Stop:     networking
# Default-Start:
# Default-Stop:      0 1 2 3 4 5 6
# Short-Description: Daemon for AWSLogs agent.
# Description:       This script is responsible for running AWSLogs agent
#                    as a daemon.
### END INIT INFO

# Version: 1.1.2-rpm

NICELEVEL=+4

export AWS_DEFAULT_REGION AWS_CONFIG_FILE HOME

# Source function library.
. /etc/rc.d/init.d/functions

prog=awslogs
exec="/usr/sbin/awslogsd"
lockfile=/var/lock/subsys/awslogs
pidfile=/var/run/awslogs.pid
mutexfile=/var/lock/awslogs.mutex

start() {
    (
    flock -w 10 -x 9
    if [ $UID -ne 0 ] ; then
        echo "Use has insufficient privilege."
        exit 4
    fi
    [ -x $exec ] || exit 5
    echo -n $"Starting $prog: "
    daemon $NICELEVEL --pidfile=$pidfile --check=${prog} "nohup $exec >/dev/null 2>&1 &"
    retval=$?
    echo
    [ $retval -eq 0 ] && touch $lockfile
    ) 9>${mutexfile}
    rm -f ${mutexfile}
}

stop() {
    (
    flock -w 10 -x 9
    if [ $UID -ne 0 ] ; then
        echo "User has insufficient privilege."
        exit 4
    fi
    echo -n $"Stopping $prog: "
    if [ -n "`pidfileofproc $prog`" ]; then
        killproc $prog
    else
        failure $"Stopping $prog"
    fi
    retval=$?
    echo
    [ $retval -eq 0 ] && rm -f $lockfile
    ) 9>${mutexfile}
    rm -f ${mutexfile}
}

restart() {
    rh_status_q && stop
    start
}

rh_status() {
    status -p /var/run/awslogs.pid -b /usr/bin/python2.7 $prog
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}


case "$1" in

    start)
        rh_status_q && exit 0
        $1
        ;;
    stop)
        rh_status_q || exit 0
        $1
        ;;
    restart)
        $1
        ;;
    status)
        rh_status
        ;;
    condrestart|try-restart)
        rh_status_q || exit 0
        restart
        ;;
    *)
        echo "Usage: $0 {start|stop|status|restart|condrestart|try-restart}"
        exit 2
        ;;

esac
exit $?